# crontab
2 10 * * * /scripts/your_script.sh >> /var/log/cron.log 2>&1
